# share/lightbar.service
[Unit]
Description=Lightbar Daemon
Requires=pi-blaster.service
Requires=dbus.service
After=pi-blaster.service
After=dbus.service

[Service]
Type=dbus
BusName=org.Lightbar
EnvironmentFile=-/etc/default/lightbar
User=root
ExecStart=/usr/bin/lightbar $DAEMON_OPTS
Restart=always
TimeoutStartSec=5
TimeoutStopSec=30
RestartSec=5

[Install]
WantedBy=multi-user.target

